# supervisord.conf
[supervisord]
nodaemon=true

[program:hephaestus]
command=hephaestus -l %(ENV_LOG_LEVEL)s
autostart=true
autorestart=true
stdout_logfile=%(ENV_STD_OUT_LOGFILE)s
stderr_logfile=%(ENV_STD_ERR_LOGFILE)s